# systemd/user/gpg-agent.service.erb
[Unit]
Description=Start gpg-agent

[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon --enable-ssh-support --write-env-file=<%= @box_homedir %>/.gnupg/gpg-agent.env

[Install]
WantedBy=default.target
